# daemon/rack.mount.in
[Mount]
Type=rack
Exec=@libexecdir@/gvfsd-rack
AutoMount=false
